# sys_files/usr/lib/systemd/system/first-boot-setup.service
[Unit]
Description=First Boot Interactive User Setup
ConditionFirstBoot=yes
Before=display-manager.service greetd.service systemd-user-sessions.service
After=local-fs.target
Before=display-manager.service greetd.service systemd-user-sessions.service

[Service]
Type=oneshot
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty2

# 2. The magic directives to take control of the screen
TTYVTAlloc=yes
TTYReset=yes
TTYVHangup=yes

ExecStart=/usr/local/bin/first-boot-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
